# docker/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/app/logs

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# Virtual Display
[program:xvfb]
command=Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
user=root
autostart=true
autorestart=true
stdout_logfile=/app/logs/xvfb.log
stderr_logfile=/app/logs/xvfb.log
environment=DISPLAY=:99

# VNC Server
[program:vnc]
command=x11vnc -display :99 -nopw -listen localhost -xkb -ncache 10 -ncache_cr -forever
user=root
autostart=true
autorestart=true
stdout_logfile=/app/logs/vnc.log
stderr_logfile=/app/logs/vnc.log
environment=DISPLAY=:99

# Redis Server
[program:redis]
command=redis-server --port 6379 --daemonize no --appendonly yes --maxmemory 512mb --maxmemory-policy allkeys-lru
user=root
autostart=true
autorestart=true
stdout_logfile=/app/logs/redis.log
stderr_logfile=/app/logs/redis.log

# Main Application
[program:autobuy-stack]
command=python src/core/orchestrator.py --start
directory=/app
user=appuser
autostart=true
autorestart=true
stdout_logfile=/app/logs/autobuy-stack.log
stderr_logfile=/app/logs/autobuy-stack.log
environment=PYTHONPATH=/app,DISPLAY=:99

# Prometheus Node Exporter
[program:node-exporter]
command=python -c "
from prometheus_client import start_http_server, Counter, Histogram, Gauge
import time
import psutil

# Create metrics
REQUEST_COUNT = Counter('autobuy_requests_total', 'Total requests')
REQUEST_LATENCY = Histogram('autobuy_request_duration_seconds', 'Request latency')
ACTIVE_TASKS = Gauge('autobuy_active_tasks', 'Active tasks')
CPU_USAGE = Gauge('system_cpu_usage_percent', 'CPU usage percentage')
MEMORY_USAGE = Gauge('system_memory_usage_percent', 'Memory usage percentage')

# Start metrics server
start_http_server(8000)
print('Metrics server started on port 8000')

# Update system metrics
while True:
    try:
        CPU_USAGE.set(psutil.cpu_percent(interval=1))
        MEMORY_USAGE.set(psutil.virtual_memory().percent)
        time.sleep(10)
    except Exception as e:
        print(f'Metrics update failed: {e}')
        time.sleep(10)
"
user=appuser
autostart=true
autorestart=true
stdout_logfile=/app/logs/metrics.log
stderr_logfile=/app/logs/metrics.log

# Health Monitor
[program:health-monitor]
command=python -c "
import asyncio
import aiohttp
import json
import time
from datetime import datetime

async def health_check():
    async with aiohttp.ClientSession() as session:
        try:
            async with session.get('http://localhost:8080/health') as response:
                if response.status == 200:
                    data = await response.json()
                    print(f'Health check OK: {json.dumps(data)}')
                else:
                    print(f'Health check failed: {response.status}')
        except Exception as e:
            print(f'Health check error: {e}')

async def main():
    while True:
        await health_check()
        await asyncio.sleep(30)

asyncio.run(main())
"
user=appuser
autostart=true
autorestart=true
stdout_logfile=/app/logs/health-monitor.log
stderr_logfile=/app/logs/health-monitor.log

# Log Rotator
[program:log-rotator]
command=python -c "
import os
import time
import shutil
from datetime import datetime, timedelta

def rotate_logs():
    log_dir = '/app/logs'
    max_size = 100 * 1024 * 1024  # 100MB
    max_age = 7  # days

    for filename in os.listdir(log_dir):
        if filename.endswith('.log'):
            filepath = os.path.join(log_dir, filename)

            # Check file size
            if os.path.getsize(filepath) > max_size:
                timestamp = datetime.now().strftime('%Y%m%d_%H%M%S')
                backup_name = f'{filename}.{timestamp}'
                shutil.move(filepath, os.path.join(log_dir, backup_name))
                print(f'Rotated log: {filename} -> {backup_name}')

            # Clean old backups
            backup_files = [f for f in os.listdir(log_dir) if f.startswith(filename + '.')]
            for backup in backup_files:
                backup_path = os.path.join(log_dir, backup)
                if os.path.getctime(backup_path) < (datetime.now() - timedelta(days=max_age)).timestamp():
                    os.remove(backup_path)
                    print(f'Removed old backup: {backup}')

while True:
    try:
        rotate_logs()
        time.sleep(3600)  # Run every hour
    except Exception as e:
        print(f'Log rotation failed: {e}')
        time.sleep(3600)
"
user=appuser
autostart=true
autorestart=true
stdout_logfile=/app/logs/log-rotator.log
stderr_logfile=/app/logs/log-rotator.log
